# deploy/systemd/kariya-sns-ui.service.d/20-envelope-credentials.conf.example
# Example only. Install as a root-owned 0644 systemd drop-in after the
# encrypted credential blobs have been provisioned with root-owned 0600 modes.
# This file contains identifiers only; it contains no key material.
[Service]
LoadCredentialEncrypted=ksns-transaction-kek-current
# During a bounded rotation overlap, also enable:
# LoadCredentialEncrypted=ksns-transaction-kek-previous

UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
